# alsa-utils-1.0.25/alsactl/90-alsa-restore.rules.in
ACTION=="add", SUBSYSTEM=="sound", KERNEL=="controlC*", KERNELS=="card*", \
        TEST=="@sbindir@/alsactl", RUN+="@sbindir@/alsactl restore $attr{number}"
